# debian/cron.d
# vim: ft=crontab

* * * * *   root   ubic-watchdog ubic.watchdog    >>/var/log/ubic/cron-watchdog.log  2>>/var/log/ubic/cron-watchdog.err.log
