# debian/postrm
#!/bin/sh
set -e

if [ \( "$1" = "purge" \) -o \( "$1" = "remove" \) ] ; then
    update-rc.d ubic.watchdog remove >/dev/null || exit $?
fi

#DEBHELPER#
