# debian/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ] ; then
    update-rc.d ubic-ping remove >/dev/null || exit $?
fi

#DEBHELPER#

