# debian/postinst
#!/bin/sh
set -e

for dir in /var/lib/ubic/watchdog /var/lib/ubic/simple-daemon/pid /var/lib/ubic/lock /var/lib/ubic/ubic-daemon /var/lib/ubic/tmp; do
    chmod 777 $dir
    chmod +t $dir
done

if [ -x "/etc/init.d/ubic-ping" ]; then
    update-rc.d ubic-ping defaults >/dev/null
    invoke-rc.d ubic-ping start || exit $?
fi

